# crontab
0 * * * * bundle exec rake "sample_data"
